# dot_config/systemd/user/bing-wallpaper.service
[Unit]
Description=Update Bing wallpaper

[Service]
ExecStart=%h/.local/bin/bing-wallpaper.sh

